# init.d/urandom
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

URANDOM_SEED=${URANDOM_SEED:-/var/run/random-seed}

depend() {
	need localmount
}

start() {
	[[ -c /dev/urandom ]] || return
	if [[ -f ${URANDOM_SEED} ]] ; then
		cat ${URANDOM_SEED} > /dev/urandom
	fi
	if ! rm -f ${URANDOM_SEED} &> /dev/null ; then
		ewarn "Skipping ${URANDOM_SEED} initialization (ro root?)"
		return 0
	fi
	ebegin "Initializing random number generator"
	(
		umask 077
		dd if=/dev/urandom of=${URANDOM_SEED} count=1 &> /dev/null
	)
	eend $? "Error initializing random number generator"
}

stop() {
	[[ -n ${CDBOOT} ]] && return 0

	ebegin "Saving random seed"
	# Carry a random seed from shut-down to start-up;
	# see documentation in linux/drivers/char/random.c
	(
		umask 077
		dd if=/dev/urandom of=${URANDOM_SEED} count=1 &> /dev/null
	)
	eend $? "Failed to save random seed"
}


# vim:ts=4
